# apache.ini
[program:apache]
command = apachectl -DFOREGROUND
