# fpm.ini
[program:fpm]
command = /opt/remi/php71/root/usr/sbin/php-fpm -F
