# mysql.ini
[program:mysql]
command = mysqld --user=mysql
