# ssh.ini
[program:ssh]
command = /usr/sbin/sshd -D -p 2222
